# pubclients.conf
description "Check pub for wifi clients"
author "er0k"

start on startup
stop on shutdown
respawn

script
  sudo -u root /home/er0k/bin/pubClients
end script
